# build/packages/portoshim/debian/portoshim.service
[Unit]
Description=CRI plugin for Porto
After=network-online.target
Requires=network-online.target
ConditionVirtualization=!container
Documentation=https://github.com/ten-nancy/portoshim

[Service]
ExecStart=/usr/sbin/portoshim
Restart=on-failure
KillSignal=SIGTERM
KillMode=process
TimeoutStopSec=360
TimeoutStartSec=360
Delegate=true

[Install]
WantedBy=multi-user.target
